# glewlwyd.service
[Unit]
Description=Glewlwyd OAuth2 authentication provider
After=network.target

[Service]
Type=simple
EnvironmentFile=-/usr/local/etc/glewlwyd/glewlwyd.conf
ExecStart=/usr/local/bin/glewlwyd --config-file=/usr/local/etc/glewlwyd/glewlwyd.conf
User=pi
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=glewlwyd
